# minimize_testcase/test_fw/squashfs-root/usr/lib/opkg/info/dnsrecorder.list
/usr/bin/dnsrecorder
/etc/init.d/dnsrecorder
