# minimize_testcase/test_fw/squashfs-root/usr/lib/opkg/info/gemtek-base-files.list
/etc/init.d/alg
/bin/debug_system.sh
/usr/bin/macClone.sh
/bin/wlanbridge_test.sh
/usr/bin/easymesh_wps.sh
/etc/hotplug.d/iface/16-wancheck
/lib/hummer/ipv6.sh
/usr/bin/timezone.sh
/lib/hummer/network.sh
/usr/bin/sync_guest_client.sh
/bin/wireless.sh
/etc/hotplug.d/client/00-wireless-client
/usr/bin/wan.sh
/lib/hummer/api.sh
/usr/bin/checkFW.sh
/bin/qos_speedtest.sh
/usr/bin/freemem.sh
/usr/bin/systemreboot.sh
/usr/bin/ddns.sh
/etc/init.d/ddns
/bin/apcli_check.sh
/bin/apcli_monitor.sh
/usr/bin/miniupnpd.sh
/bin/dhcpInfo.sh
/lib/hummer/network/wan.sh
/usr/bin/ledstatus.sh
/etc/init.d/nbnsd
/usr/bin/EasyMesh_openwrt_gemtek.sh
/usr/bin/easymesh_led.sh
/lib/hummer/system.sh
/etc/syslog.conf
/usr/bin/checkinternet.sh
/etc/init.d/autofw
/lib/hummer/qos.sh
/usr/bin/checkpw.sh
/etc/version
/bin/debug_mpt.sh
/www/json/wlanMPTable.json
/bin/wps_action.sh
/lib/hummer/network/lan.sh
/usr/bin/testeasymesh.sh
/usr/bin/vlan_setup.sh
/usr/bin/checklanoverlay.sh
/lib/hummer/network/routing.sh
/etc/init.d/logrotate
/bin/get_hosts.sh
/lib/functions/gemtek.sh
/usr/bin/easymesh_main.sh
/usr/bin/updatesysconfig.sh
/etc/init.d/syslog
/usr/bin/easymesh.sh
/bin/ipcalc.sh
/lib/hummer/state.sh
/www/cgi-bin/debugmsg.cgi
