# minimize_testcase/test_fw/squashfs-root/usr/lib/opkg/info/hwnat.list
/etc/init.d/hwnat
/bin/hwnat
/etc/config/hwnat
