# minimize_testcase/test_fw/squashfs-root/usr/lib/opkg/info/mdnsd.list
/usr/lib/libdns_sd.so.1
/usr/lib/libdns_sd.so
/usr/sbin/mdnsd
/etc/init.d/mdnsd
