# minimize_testcase/test_fw/squashfs-root/usr/lib/opkg/info/nodogsplash.list
/etc/nodogsplash/htdocs/images/normal_left.gif
/etc/config/nodogsplash
/usr/bin/nodogsplash
/lib/upgrade/keep.d/nodogsplash
/usr/bin/ndsctl
/etc/nodogsplash/htdocs/images/Background.jpg
/etc/nodogsplash/htdocs/images/normal_right.gif
/etc/nodogsplash/htdocs/images/top_to_bottom_gradient.jpg
/etc/nodogsplash/htdocs/images/normal_middle.gif
/etc/init.d/nodogsplash
/etc/nodogsplash/htdocs/splash.html
/etc/nodogsplash/htdocs/images/guest.css
/etc/nodogsplash/htdocs/images/splash.jpg
/etc/nodogsplash/htdocs/infoskel.html
/etc/nodogsplash/htdocs/images/Divider.jpg
