# minimize_testcase/test_fw/squashfs-root/usr/lib/opkg/info/quagga.list
/usr/sbin/quagga.init
/etc/init.d/quagga
