# minimize_testcase/test_fw/squashfs-root/usr/lib/opkg/info/urlfilter.list
/etc/init.d/urlfilter
/usr/bin/urlfilter
